# samples/django-celery-app/config/sample-datastore/supervisor.d/postgresql.ini
[program:postgresql]
command=/usr/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql_err.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
environment=PGDATA="/var/lib/postgresql/data"
stopsignal=INT
stopwaitsecs=30