# samples/django-celery-app/config/sample-datastore/supervisor.d/redis.ini
[program:redis]
command=/usr/bin/redis-server /etc/redis.conf
user=redis
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/redis/redis.log
stderr_logfile=/var/log/redis/redis_err.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stopsignal=QUIT
stopwaitsecs=10