# srv/lxc-compose/scripts/internal/lxc-port-forward.service
[Unit]
Description=LXC Compose Port Forwarding Rules
After=network.target lxc-net.service lxd.service
Wants=lxc-net.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/lxc-compose port apply
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target