# samples/django-minimal/config/supervisor.d/postgresql.ini
[program:postgresql]
command=/usr/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql_err.log
environment=PGDATA="/var/lib/postgresql/data"